# etc/hkd.service
[Unit]
Description=Hot Keys Daemon

[Service]
Type=simple
ExecStart=/usr/bin/hkd
Restart=always

[Install]
WantedBy=multi-user.target
